# {{cookiecutter.project_name}}/config/{{cookiecutter.project_name}}.supervisor.ini
[program:{{ cookiecutter.project_name }}]
directory={{ cookiecutter.app_path }}/{{ cookiecutter.project_name }}
command=gunicorn -b unix:/tmp/gunicorn.{{ cookiecutter.project_name }}.sock app.app:app
user={{ cookiecutter.app_user }}
